# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:cmd]
command=env | grep _ >> /etc/environment && /usr/local/bin/run